# TestCases/scripts/protheus.conf
# protheus.conf
start on filesystem
exec python3 "/home/dev/Protheus/TestCases/scripts/service.py"
#exec python3 "path_to_protheus/TestCases/scripts/service.py"
